# deploy/systemd/media-pro.service
[Unit]
Description=media-pro FastAPI gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/nax/media-pro
EnvironmentFile=/etc/media-pro/media-pro.env
ExecStart=/home/nax/media-pro/.venv/bin/uvicorn gateway.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
